# crontab/django
# Used for django_cron
*/1 * * * * source /home/clapclapp/.bashrc && source /home/clapclapp/venv/wagtail/bin/activate && python /home/clapclapp/wagtail/fech/manage.py runcrons > /home/clapclapp/cronjob.log