# crontab/docker
# Used for django_cron
*/1 * * * * /usr/bin/docker exec fech_backend python manage.py runcrons > /home/ubuntu/cronjob.log